# scripts/prune-containers.timer
[Unit]
Description=Weekly Podman cleanup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
