# profile/root/etc/systemd/system/geckoforge-ssh-hardening.service
[Unit]
Description=Geckoforge first boot SSH hardening
ConditionFirstBoot=yes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/firstboot-ssh-hardening.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
